# plugin/hz-gate/hz-gate.conf
[program:hz_gate]
command=/bin/bash -c "if [ ! -d "./hz-gate" ]; then tar -zxvf hz-gate.tar.gz; fi; cd hz-gate;chmod +x ./run.sh ./root/start-gate.sh ./root/DGPServer ./root/HuaZhi;./run.sh"
directory=/data/conf.d
user=root
autostart=true
autorestart=false
stopasgroup=true
killasgroup=true
stdout_logfile=/data/log/hz_gate.log
redirect_stderr=true
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
