# plugin/ls-gate/ls-gate.conf
[program:ls_gate]
command=/bin/bash -c "if [ ! -d "./ls-gate" ]; then tar -zxvf ls-gate.tar.gz; fi; cd ls-gate;chmod +x ./run.sh;./run.sh"
directory=/data/conf.d
user=root
autostart=true
autorestart=false
stopasgroup=true
killasgroup=true
stdout_logfile=/data/log/ls_gate.log
redirect_stderr=true
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
